# src/pold.service.in
[Unit]
Description=GENIVI Policy Daemon

[Service]
Type=dbus
BusName=de.bmw.net.pold
ExecStart=@prefix@/bin/pold
